# build/python_orocos_kdl/devel/etc/catkin/profile.d/python_orocos_kdl_site_packages.tcsh
export PYTHONPATH=/home/ur3/catkin_FrankaPanda/devel/lib/python3/dist-packages:$PYTHONPATH
